# debian/rsbackup.conffiles
/etc/rsbackup/config
/etc/rsbackup/defaults
/etc/cron.hourly/rsbackup
/etc/cron.daily/rsbackup
/etc/cron.weekly/rsbackup
